# niri-gamemode-daemon/niri-gamemode-daemon.service
[Unit]
Description=Niri GameMode Daemon - Auto-toggle keyboard game mode for gamescope
After=graphical-session.target
Documentation=https://github.com/yourusername/niri-gamemode-daemon

[Service]
Type=simple
ExecStart=%h/.cargo/bin/niri-gamemode-daemon
Restart=always
RestartSec=5

# Environment variables (uncomment and set if you want to specify a specific keyboard)
# Environment="KEYBOARD_VID=0x362D"
# Environment="KEYBOARD_PID=0x0210"

# Give the service access to the user's session bus and display
Environment="DISPLAY=:0"

# Reduce logging verbosity in production (comment out for debugging)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
